# scripts/auth-gate.service
[Unit]
Description=Auth Gate API Gateway
After=network.target

[Service]
Type=simple
User=auth-gate
Group=auth-gate
WorkingDirectory=/opt/auth-gate
ExecStart=/opt/auth-gate/auth-gate
Restart=always
RestartSec=5

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/auth-gate

[Install]
WantedBy=multi-user.target
